# systemd/schoolLib.service.example
[Unit]
Description=Run the All Saints Simple School Lirbary
After=network.target

[Service]
ExecStart=/home/<aUser>/.local/bin/schoolLib

[Install]
WantedBy=default.target
